# files/supervisor/btcd.conf
[program:btcd]
command=/home/ubuntu/go/bin/btcd --configfile=/home/ubuntu/.btcd/btcd.conf

user=ubuntu

startsecs=3
stopwaitsecs=10
autostart=true
autorestart=true
stopsignal=TERM
